# Linux-server-setup/docker-security-scan.timer
[Unit]
Description=Daily Docker Security Scanner Timer
After=docker.service

[Timer]
OnCalendar=*-*-* 03:00:00
RandomizedDelaySec=1800
Persistent=true

[Install]
WantedBy=timers.target
